# gatus-consul/gatus/s6-rc.d/gatus/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -x

export GATUS_DATA_DIR="${GATUS_DATA_DIR:-"/gatus/data"}"

install -v -d -m 755 "$GATUS_DATA_DIR"

exec gatus
