# git-nginx/s6-rc.d/nginx/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -e
set -x

exec nginx -g 'daemon off;'
