# git-nginx/s6-rc.d/nginx/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -e
set -x

# Ensure nginx user ID is updated
groupmod -o -g "$PGID" nginx
usermod -o -u "$PUID" nginx

exec nginx -g 'daemon off;'
